# packages/sys-apps/razercfg/files/systemd/razerd.service
[Unit]
Description=Advanced configuration of Razer mice (DeathAdder, Krait, Lachesis)

[Service]
ExecStart=/usr/sbin/razerd

[Install]
WantedBy=multi-user.target
